# ansible/roles/pico-web/templates/ctf-daemon.service.j2
[Unit]
Description= ctf daemon runner

[Service]
Type=simple
Environment="APP_SETTINGS_FILE={{ web_config_dir }}/deploy_settings.py"
ExecStart=/usr/local/bin/daemon_manager -d {{ daemon_src_dir }} cache_stats share_instances

[Install]
WantedBy=multi-user.target
